# quick-auth.service.example
# ── Quick Auth — systemd 服务配置 ──
# 安装方法：
#   1. 将本文件复制到 /etc/systemd/system/quick-auth.service
#   2. 修改 WorkingDirectory、User、Group 等路径为实际值
#   3. 执行以下命令：
#      sudo systemctl daemon-reload
#      sudo systemctl enable quick-auth
#      sudo systemctl start quick-auth
#   4. 查看状态：
#      sudo systemctl status quick-auth
#   5. 查看日志：
#      sudo journalctl -u quick-auth -f

[Unit]
Description=Quick Auth - Email verification based authentication service
After=network.target
Wants=network-online.target

[Service]
Type=exec

# ── 工作目录（项目根目录，确保 .env 文件在此路径下） ──
WorkingDirectory=/opt/quick-auth

# ── 运行用户和组 ──
User=node
Group=node

# ── 启动命令（使用 tsx 直接运行 TypeScript） ──
ExecStart=/usr/bin/node /usr/bin/tsx /opt/quick-auth/src/index.ts

# ── 重启策略 ──
Restart=always
RestartSec=5

# ── 环境变量（也可通过项目根目录的 .env 文件加载） ──
Environment=NODE_ENV=production
Environment=PORT=3000

# ── 日志 ──
StandardOutput=journal
StandardError=journal
SyslogIdentifier=quick-auth

# ── 安全与资源限制 ──
NoNewPrivileges=true
ProtectHome=true
ProtectSystem=full
PrivateTmp=true
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
